# hetzner/config/audit.rules
# Clean up existing rules
-D

# Buffer size (increase if audit backlog warnings appear)
-b 8192

# Fail gracefully (1 = printk, 2 = panic)
-f 1

# --- Authentication & Authorization ---
-w /etc/passwd -p wa -k auth
-w /etc/shadow -p wa -k auth
-w /etc/group -p wa -k auth
-w /etc/gshadow -p wa -k auth
-w /etc/sudoers -p wa -k auth
-w /etc/sudoers.d/ -p wa -k auth
-w /etc/pam.d/ -p wa -k auth
-w /var/log/auth.log -p ra -k auth
-w /var/log/faillog -p ra -k auth

# --- User/Group Management ---
-w /usr/sbin/useradd -p x -k user-mgmt
-w /usr/sbin/userdel -p x -k user-mgmt
-w /usr/sbin/usermod -p x -k user-mgmt
-w /usr/sbin/groupadd -p x -k user-mgmt
-w /usr/sbin/groupdel -p x -k user-mgmt
-w /usr/sbin/groupmod -p x -k user-mgmt

# --- Sudo & Privilege Escalation ---
-w /usr/bin/sudo -p x -k privilege
-w /usr/bin/su -p x -k privilege

# --- Cron ---
-w /etc/crontab -p wa -k cron
-w /etc/cron.d/ -p wa -k cron
-w /etc/cron.daily/ -p wa -k cron
-w /etc/cron.weekly/ -p wa -k cron
-w /etc/cron.monthly/ -p wa -k cron
-w /var/spool/cron/ -p wa -k cron

# --- System Services ---
-w /etc/systemd/ -p wa -k systemd
-w /usr/lib/systemd/ -p wa -k systemd
-w /usr/bin/systemctl -p x -k systemd

# --- Network ---
-w /etc/ufw/ -p wa -k network
-w /usr/sbin/ufw -p x -k network
-w /etc/hosts -p wa -k network
-w /etc/hostname -p wa -k network
-w /etc/resolv.conf -p wa -k network

# --- SSH & Tailscale ---
-w /etc/ssh/ -p wa -k ssh
-w /root/.ssh/ -p wa -k ssh
-w /usr/bin/tailscale -p x -k tailscale

# --- Package Management ---
-w /usr/bin/apt -p x -k packages
-w /usr/bin/apt-get -p x -k packages
-w /usr/bin/dpkg -p x -k packages

# --- Docker ---
-w /usr/bin/dockerd -k docker
-w /usr/bin/docker -k docker
-w /etc/docker/ -p wa -k docker
-w /var/lib/docker -p wa -k docker
-w /usr/lib/systemd/system/docker.service -p wa -k docker

# --- Kernel & Boot ---
-w /etc/sysctl.conf -p wa -k kernel
-w /etc/sysctl.d/ -p wa -k kernel
-w /boot/ -p wa -k boot

# Make rules immutable (requires reboot to change)
-e 2
